# infra/systemd/agentops-relay.service
[Unit]
Description=Agent Ops Relay Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/clawbot.env
WorkingDirectory=/opt/agent-ops
ExecStart=/usr/bin/python3 /opt/agent-ops/scripts/relay_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target