# infra/systemd/clawbot-monitor.timer
[Unit]
Description=Run Clawbot Monitor every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Unit=clawbot-monitor.service

[Install]
WantedBy=timers.target